# priv/templates/smartos/manifest18.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="{{package_install_name}}">
  <service name="application/{{package_install_name}}" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name="{{package_install_name}}-epmd" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/network/{{package_install_name}}-epmd:default" />
    </dependency>
    <method_context working_directory="/tmp" project="{{package_install_name}}">
      <method_credential user="{{package_install_user}}" group="{{package_install_group}}" />
      <method_environment>
        <envvar name="HOME" value="{{platform_base_dir}}" />
        <envvar name="LOGNAME" value="{{package_install_user}}" />
        <envvar name="PATH" value="/opt/local/gnu/bin:/opt/local/bin:/opt/local/sbin:/usr/bin:/usr/sbin" />
        <envvar name="LD_PRELOAD_32" value="/lib/libumem.so.1" />
        <envvar name="LD_PRELOAD_64" value="/lib/amd64/libumem.so.1:/opt/local/gcc47/lib/amd64/libgcc_s.so.1"/>
        <envvar name="UMEM_OPTIONS" value="allocator=best" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="{{platform_bin_dir}}/{{package_install_name}} start" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="{{platform_bin_dir}}/{{package_install_name}} stop" timeout_seconds="60" />
    <exec_method type="method" name="restart" exec="{{platform_bin_dir}}/{{package_install_name}} restart" timeout_seconds="60" />
    <stability value="Stable" />
    <template>
      <common_name>
        <loctext xml:lang="C">{{package_shortdesc}}</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
